# cloudflare-dyndns@.timer
[Unit]
Description=Check Cloudflare DNS periodically (%i)

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
